# config/autostart/autostart.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Custom Autostart Programs
Exec=myautostart
StartupNotify=false
Type=Application
