# config/autostart/autostart.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Custom Autostart Programs
Exec=bash -c $(xdg-user-dir)/bin/myautostart
StartupNotify=false
Type=Application
